# templates/etc/sv/spawn-fcgi/run
#!/bin/sh
exec /usr/bin/spawn-fcgi -n -u {{_cgi_user}} -s /run/spawn-fcgi.sock -- /usr/sbin/fcgiwrap
