# templates/etc/sv/git-daemon/run
#!/bin/sh

exec chpst -u nobody:nogroup \
  /usr/libexec/git-core/git-daemon --syslog --base-path=/var/lib/git/public
